# backup-full-2025-11-18-09-01-22/pingit.service
[Unit]
Description=PingIT - Network Monitoring Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as root to allow ICMP socket creation (required for ping)
# Consider using setcap as alternative if security is a concern
User=root
Group=root
WorkingDirectory=/opt/pingit
ExecStart=/usr/bin/python3 /opt/pingit/pingit.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pingit

[Install]
WantedBy=multi-user.target
