# backup-full-2025-11-15-18-51-46/pingit.service
[Unit]
Description=PingIT - Network Monitoring Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pingit
Group=pingit
WorkingDirectory=/opt/pingit
ExecStart=/usr/bin/python3 /opt/pingit/pingit.py --config /etc/pingit/config.yaml --log-dir /var/log/pingit
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pingit

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/pingit /var/lib/pingit

[Install]
WantedBy=multi-user.target
